# hosts-updater.service
# hosts-updater.service

[Unit] 
Description= hosts-updater

[Service] 
Type= simple 

ExecStart= /usr/local/bin/hosts-updater.sh

